# root/etc/s6-overlay/s6-rc.d/svc-claude-code-ui/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Run Claude Code UI as a long-running service (as abc user)
# Depends on svc-claude-code-startup which installs Claude Code CLI
# Claude Code UI is built from source at /opt/claude-code-ui during image build
# Directory /config/.claude-code-ui is created by init-claude-code-config (as root, with correct ownership)

export HOME="/config"
export PATH="/config/.npm-global/bin:/config/.local/bin:$PATH"
export CLAUDE_CLI_PATH="/config/.local/bin/claude"
export DATABASE_PATH="/config/.claude-code-ui/auth.db"
export PORT=3001

cd /opt/claude-code-ui
exec s6-setuidgid abc node server/index.js
