# root/etc/s6-overlay/s6-rc.d/svc-claude-code-ui/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Run Claude Code UI as a long-running service (as abc user)
# Depends on svc-claude-code-startup which installs the npm package

# Set up environment for the abc user
export HOME="/config"
export PATH="/config/.npm-global/bin:/config/.local/bin:$PATH"
export CLAUDE_CLI_PATH="/config/.local/bin/claude"
export DATABASE_PATH="/config/.claude-code-ui/auth.db"
export PORT=3001

# Ensure database directory exists
mkdir -p /config/.claude-code-ui

exec s6-setuidgid abc claude-code-ui
