# root/etc/s6-overlay/s6-rc.d/svc-claude-code-startup/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This runs as the abc user
exec s6-setuidgid abc /config/scripts/startup.sh