# root/etc/s6-overlay/s6-rc.d/svc-agent-os/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Run Agent-OS as a longrun service on port 3011
# This provides a mobile-first web UI for managing AI coding sessions

cd /config || exit 1

# Export PATH to include npm global binaries
export PATH="/config/.npm-global/bin:$PATH"
export HOME="/config"
export AGENT_OS_HOME="/config/.agent-os"
export AGENT_OS_PORT="3011"

# Wait for Agent-OS to be installed (startup script runs first)
until [ -x /config/.npm-global/bin/agent-os ]; do
    echo "Waiting for Agent-OS to be installed..."
    sleep 5
done

# Additional wait for the initial setup to complete
if [ ! -d /config/.agent-os/repo ]; then
    echo "Waiting for Agent-OS initial setup to complete..."
    sleep 10
fi

# Start Agent-OS server in foreground mode (required for s6 longrun services)
# Note: EBADF errors are expected when running without a controlling terminal
# These are non-fatal and don't affect the web UI functionality
echo "Starting Agent-OS web server on port 3011..."
echo "Note: Terminal-related errors (EBADF) are expected and can be safely ignored"
exec s6-setuidgid abc agent-os start-foreground