# root/etc/s6-overlay/s6-rc.d/svc-agent-os/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Run Agent-OS as a longrun service on port 3011
# This provides a mobile-first web UI for managing AI coding sessions

cd /config || exit 1

# Export PATH to include npm global binaries
export PATH="/config/.npm-global/bin:$PATH"
export HOME="/config"
export USER="abc"
export AGENT_OS_HOME="/config/.agent-os"
export AGENT_OS_PORT="3011"
# Set shell to bash (container doesn't have zsh)
export SHELL="/bin/bash"

# Wait for Agent-OS to be installed (startup script runs first)
until [ -x /config/.npm-global/bin/agent-os ]; do
    echo "Waiting for Agent-OS to be installed..."
    sleep 5
done

# Additional wait for the initial setup to complete
if [ ! -d /config/.agent-os/repo ]; then
    echo "Waiting for Agent-OS initial setup to complete..."
    sleep 10
fi

# Start Agent-OS server in foreground mode
echo "Starting Agent-OS web server on port 3011..."

# Note: Terminal sessions may have issues when running as a service
# Agent-OS works best when run interactively from a terminal
# The web UI will work, but terminal connections may fail with EBADF errors
exec s6-setuidgid abc agent-os start-foreground