# root/etc/s6-overlay/s6-rc.d/init-claude-code-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This script runs as root during container initialization
# It ensures directories exist with proper ownership based on PUID/PGID

echo "Initializing claude-code specific directories..."

# Create application-specific directories with proper ownership
# We need to explicitly set ownership since this runs after base image's init
mkdir -p /config/scripts
mkdir -p /config/.claude
mkdir -p /config/.codex
mkdir -p /config/.gemini
mkdir -p /config/workspace
mkdir -p /config/.npm
mkdir -p /config/.npm-global

# Set ownership on directories we created using PUID/PGID from environment
echo "Setting ownership on AI assistant directories (PUID=${PUID}, PGID=${PGID})..."
chown -R "${PUID}:${PGID}" /config/.claude
chown -R "${PUID}:${PGID}" /config/.codex
chown -R "${PUID}:${PGID}" /config/.gemini
chown -R "${PUID}:${PGID}" /config/.npm
chown -R "${PUID}:${PGID}" /config/.npm-global
chown -R "${PUID}:${PGID}" /config/scripts
chown -R "${PUID}:${PGID}" /config/workspace

# No longer copying startup script to /config - it runs directly from /defaults
# This ensures every container restart uses the fresh version from the image

# Fix permissions on docker socket if mounted
if [[ -S /var/run/docker.sock ]]; then
    echo "Docker socket detected. Configuring access for abc user..."
    DOCKER_GID=$(stat -c '%g' /var/run/docker.sock)

    if [ "${DOCKER_GID}" = "0" ]; then
        # Socket owned by root group (common on Unraid) - create a docker group
        # and fix socket permissions so non-root users can access it
        if ! getent group docker >/dev/null 2>&1; then
            groupadd docker
        fi
        chgrp docker /var/run/docker.sock
        chmod g+rw /var/run/docker.sock
    else
        # Ensure docker group exists with the correct GID matching the socket
        if ! getent group docker >/dev/null 2>&1; then
            groupadd -g "${DOCKER_GID}" docker
        elif [ "$(getent group docker | cut -d: -f3)" != "${DOCKER_GID}" ]; then
            # docker group exists but with wrong GID (e.g. from docker-ce-cli package)
            groupmod -g "${DOCKER_GID}" docker
        fi
    fi

    # Add abc user to docker group
    usermod -aG docker abc
fi

echo "Claude-code directories created with proper ownership."