# root/etc/s6-overlay/s6-rc.d/init-claude-code-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# This script runs as root during container initialization
# It ensures directories exist with proper ownership based on PUID/PGID

echo "Initializing claude-code specific directories..."

# Only create application-specific directories
# Let LinuxServer base image handle overall /config ownership
# We just ensure our specific directories exist - ownership will be handled by base image
mkdir -p /config/scripts
mkdir -p /config/.claude
mkdir -p /config/workspace
mkdir -p /config/.npm

# Copy startup script if it doesn't exist
# Base image will handle ownership via lsiown
if [[ ! -f /config/scripts/startup.sh ]]; then
    cp /defaults/startup.sh /config/scripts/startup.sh
    chmod +x /config/scripts/startup.sh
fi

# Fix permissions on docker socket if mounted
if [[ -S /var/run/docker.sock ]]; then
    echo "Docker socket detected. Adding user to docker group..."
    # Get the GID of the docker socket
    DOCKER_GID=$(stat -c '%g' /var/run/docker.sock)
    
    # Create docker group with the same GID if it doesn't exist
    if ! getent group docker >/dev/null 2>&1; then
        groupadd -g "${DOCKER_GID}" docker
    fi
    
    # Add abc user to docker group
    usermod -aG docker abc
fi

echo "Claude-code directories created. Base image will handle ownership."